# systemd/LG_Buddy_screen.service
[Unit]
Description=LG Buddy Screen Monitor Service
After=graphical-session.target
PartOf=graphical-session.target

[Service]
ExecStart=/usr/bin/lg-buddy monitor
Restart=always
RestartSec=10

[Install]
WantedBy=graphical-session.target
